# systemd-01/logwatcher.timer
[Unit]
Description=Run a logwatch service every 30 seconds

[Timer]
# Run every 30 secs
OnUnitActiveSec=30

[Install]
WantedBy=multi-user.target

